# sshuttle-gui/etc/systemd/system/sstp-connector.service
[Unit]
Description=SSTP VPN Tunneling Application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/etc/sstp-connector/connect-systemd.sh
ExecStop=/etc/sstp-connector/stop-connect.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
