# sshuttle-gui/all/etc/systemd/system/sshuttle.service
[Unit]
Description=sshuttle
After=network.target

[Service]
Type=simple
ExecStart=/etc/sshuttle-gui/connect.sh

[Install]
WantedBy=multi-user.target